# qtee_supplicant/qteesupplicant.service.in
# Copyright (c) Qualcomm Technologies, Inc. and/or its subsidiaries.
# SPDX-License-Identifier: BSD-3-Clause-Clear
[Unit]
Description=QTEE Supplicant Service
BindsTo=dev-tee0.device
After=dev-tee0.device sfsconfig.service

[Service]
Type=exec
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/qtee_supplicant
Restart=always

[Install]
WantedBy=multi-user.target
